# deploy/supervisor_profiles_api.conf
[program:profiles_api]
environment =
  DEBUG=0
  SECRET_KEY=django-insecure-feewdetz2z+15()(pn)rj@ws68s!l%)a+mp#q!2gv1x)80*w*i
command = /usr/local/apps/profiles-rest-api/env/bin/uwsgi --http :9000 --wsgi-file /usr/local/apps/profiles-rest-api/profiles_project/wsgi.py
directory = /usr/local/apps/profiles-rest-api/
user = root
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/profiles_api.log
stderr_logfile = /var/log/supervisor/profiles_api_err.log
